# scripts/musicai.service
[Unit]
Description=Magic AI Music Box - AI-Powered Auto-Tune System
Documentation=https://github.com/ruvnet/musicai
After=network.target sound.target pulseaudio.service
Wants=network.target sound.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/musicai

# Main command
ExecStart=/usr/bin/npm start

# Pre-start checks
ExecStartPre=/bin/sh -c 'test -d /home/pi/musicai'
ExecStartPre=/bin/sh -c 'test -f /home/pi/musicai/package.json'

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Environment
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=1024
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Resource limits
MemoryLimit=1G
MemoryHigh=800M
CPUQuota=80%
TasksMax=100

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/musicai/data /home/pi/musicai/logs

# Audio access
SupplementaryGroups=audio

# Standard output and error logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=musicai

[Install]
WantedBy=multi-user.target
